# backup/crontab
0 * * * * /usr/local/bin/backup.sh >> /backups/cron.log 2>&1
